# target/linux/layerscape/base-files/etc/init.d/hclink_device
#!/bin/sh /etc/rc.common

START=10
STOP=98

#add by liji to modify ethx order: ethx->ifacex
set_interface_name(){
	flag="$(cat /proc/cmdline | grep hcen104)"
	if [ "$flag" ]; then
	    ip link set eth3 name iface0
	    ip link set eth0 name iface1
	    ip link set eth4 name iface2
	    ip link set eth2 name iface3
	    ip link set eth1 name iface4
	fi
	flag="$(cat /proc/cmdline | grep hcen1002)"
	if [ "$flag" ]; then
	    ip link set eth1 name iface0
	    ip link set eth2 name iface1
	    ip link set eth4 name iface2
	    ip link set eth0 name iface3
	    ip link set eth3 name iface4
	fi
	flag="$(cat /proc/cmdline | grep ag1002at)"
	if [ "$flag" ]; then
	    ip link set eth3 name iface0
	    ip link set eth0 name iface1
	    ip link set eth4 name iface2
	    ip link set eth2 name iface3
	    ip link set eth1 name iface4
	fi
}

#add by liji to set wifi reg to CN
set_iw_reg_to_cn(){
	iw reg set CN
}

start(){
	echo "hclink device start..."
	set_interface_name
	set_iw_reg_to_cn
	#hclink-liji: I found reset usually no effact,so change it mod
	chmod 777 /etc/rc.button/reset
	#boot lte
	lte up
}

stop(){
	echo "nothing to stop"
}

restart(){
	echo "please stop and start"
}